# openwrt/luci-app-photondns/root/etc/uci-defaults/40_luci-photondns
#!/bin/sh
# Refresh the LuCI dispatch cache so the new photondns menu appears.
# (Default config, init script and rule files ship with the photondns
# daemon package itself.)
rm -f /tmp/luci-indexcache* /tmp/luci-modulecache/* 2>/dev/null
exit 0
